# systemd/eeebot-guarded-evolution.service
[Unit]
Description=eeebot Guarded Self-Evolution Apply/Health/Rollback
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/ozand/herkoot/Projects/nanobot
Environment=PYTHONPATH=.
EnvironmentFile=%h/.config/eeebot-self-improving.env
ExecStart=/usr/bin/env python3 /home/ozand/herkoot/Projects/nanobot/scripts/guarded_self_evolve.py
