# ops/dashboard/systemd/nanobot-ops-dashboard-collector.service
[Unit]
Description=Nanobot Ops Dashboard Collector
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ozand/herkoot/Projects/nanobot/ops/dashboard
Environment=PYTHONPATH=src
EnvironmentFile=%h/.config/nanobot-ops-dashboard.env
ExecStart=/home/ozand/herkoot/Projects/nanobot/ops/dashboard/scripts/run_collector.sh
Restart=on-failure
RestartSec=5
MemoryMax=512M
RuntimeMaxSec=12h

[Install]
WantedBy=default.target
